# infra/lilynet/config/local.d/routing-policy.start
#!/bin/sh

{% macro table(table, priority) %}
ip -4 rule delete table {{ table }} || true
ip -4 rule add priority {{ priority }} table {{ table }}

ip -6 rule delete table {{ table }} || true
ip -6 rule add priority {{ priority }} table {{ table }}
{% endmacro %}

{{ table("babeld", 8192) }}
{{ table("bird", 16384) }}
