# infra/lilynet/config/local.d/routing-policy.start
#!/bin/sh

ip -4 rule delete table babeld || true
ip -4 rule add priority 8192 table babeld

ip -6 rule delete table babeld || true
ip -6 rule add priority 8192 table babeld

ip -4 rule delete table bird || true
ip -4 rule add priority 16384 table bird

ip -6 rule delete table bird || true
ip -6 rule add priority 16384 table bird
